# docker/supervisord.conf
# ─────────────────────────────────────────────────────────────────────────────
# Supervisord — runs Nginx + Spring Boot in a single container
# ─────────────────────────────────────────────────────────────────────────────

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:metahub-backend]
command=java %(ENV_JAVA_OPTS)s -jar /app/metahub.jar
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
